# systemd/four_in_a_row.message_consumer.service
[Unit]
Description=Four In A Row Message Consumer
After=network.target

[Service]
User=four_in_a_row
Group=four_in_a_row
Type=simple
WorkingDirectory=/opt/four_in_a_row
ExecStart=/opt/four_in_a_row/venv/bin/four_in_a_row run-message-consumer
Restart=always

[Install]
WantedBy=multi-user.target
